# steam-bigpicture.service
[Unit]
Description=Set volume for Steam Big Picture mode

[Service]
Type=simple
ExecStart=

[Install]
WantedBy=multi-user.target
